# rootfs/etc/s6-overlay/s6-rc.d/backup_api/run
#!/usr/bin/with-contenv bashio

bashio::log.info "---------Starting python API sript--------"

python3 /etc/s6-overlay/scripts/check_backup_server_status.py