# rootfs/etc/s6-overlay/s6-rc.d/backup_api/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Fleet Assistant
# Runs backup service
# ==============================================================================
#declare serverip

bashio::log.info "Running API backup.."


#if bashio::config.has_value "server.tunnelip"; then
#    serverip=$(bashio::config "server.tunnelip")
#fi

#python3 /usr/bin/check_backup_server_status.py --FleetAssistantServerIP=serverip

python3 /scripts/check_backup_server_status.py
