# rootfs/etc/s6-overlay/s6-rc.d/backup_api/run
#!/usr/bin/with-contenv bashio

bashio::log.info "---------Starting python API sript--------"

exec python3 -u /etc/s6-overlay/scripts/check_backup_server_status.py --FleetAssistantServerIP 10.170.204.1